# scripts/newsnexus12-worker-node-orchestrator-test.service
[Unit]
Description=NewsNexus12 worker-node abbreviated test orchestration trigger
Documentation=file:/home/limited_user/applications/NewsNexus12/scripts/newsnexus12-worker-node-orchestrator-test.service
Requires=newsnexus12-worker-node.service
After=network-online.target newsnexus12-worker-node.service
Wants=network-online.target

[Service]
Type=oneshot
User=limited_user
Group=limited_user
WorkingDirectory=/home/limited_user/applications/NewsNexus12/worker-node
Environment="WORKER_NODE_URL=http://127.0.0.1:8003"
ExecStart=/home/limited_user/applications/NewsNexus12/scripts/trigger-worker-node-orchestrator-test.sh
